# rfd-start.conf
[program:bluetoothd]
command=/bin/bash -c "/usr/lib/bluetooth/bluetoothd"
redirect_stderr=true

