# supervisor/rsyslog.conf
[program:rsyslog]
command=/bin/bash -c "service rsyslog start"
autostart=true
startretries=0
autorestart=false
redirect_stderr=true
exitcodes=0
priority=5
nodaemon=true

