# init.d/lighttpd
### BEGIN INIT INFO
# Provides:             lighttpd
# Required-Start:       $local_fs $remote_fs
# Required-Stop:        $local_fs $remote_fs
# Default-Start:        2 3 4 5
# Default-Stop:         0 1 6
# Short-Description:    Start lighttpd
### END INIT INFO

start() {
        echo -n "Starting lighttpd: "
        export HM_HOME=/opt/hm
        export LD_LIBRARY_PATH=$HM_HOME/lib
        start-stop-daemon -S -q -p /var/run/lighttpd.pid --exec $HM_HOME/bin/lighttpd -- -f /opt/hm/etc/lighttpd/lighttpd.conf -m /opt/hm/lib/lighttpd/
        echo "OK"


}
stop() {
        echo -n "Stopping lighttpd: "
        start-stop-daemon -K -q -p /var/run/lighttpd.pid
        echo "OK"
}
restart() {
        stop
        start
}

case "$1" in
  start)
        start
        ;;
  stop)
        stop
        ;;
  restart|reload)
        restart
        ;;
  *)
        echo "Usage: $0 {start|stop|restart}"
        exit 1
esac

exit $?

